# ping-check.service
[Unit]
Description=Check Hosts via Ping

[Service]
ExecStart=/usr/local/bin/ping-check
